# www-servers/glassfish-bin/files/glassfish-5-init
#!/sbin/openrc-run

DOMAIN=${SVCNAME/glassfish-5-/}

depend() {
    need net
}

start() {
    ebegin "Starting GlassFish application server"

    if [[ "x${DOMAIN}" == "x" ]]; then
        su -l glassfish -c "/opt/glassfish-5/bin/asadmin start-domain"
    else
        einfo "Starting domain ${DOMAIN}"
        su -l glassfish -c "/opt/glassfish-5/bin/asadmin start-domain ${DOMAIN}"
    fi
    eend $?
}

stop() {
    ebegin "Sopping GlassFish application server"

    if [[ "x${DOMAIN}" == "x" ]]; then
        su -l glassfish -c "/opt/glassfish-5/bin/asadmin stop-domain"
    else
        einfo "Stopping domain ${DOMAIN}"
        su -l glassfish -c "/opt/glassfish-5/bin/asadmin stop-domain ${DOMAIN}"
    fi
    eend $?
}
